# deploy/linux/hermes-email-bridge.service
[Unit]
Description=Hermes email bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__BRIDGE_USER__
Group=__BRIDGE_GROUP__
EnvironmentFile=__ENV_FILE__
WorkingDirectory=__WORKSPACE__
ExecStartPre=/usr/bin/sudo -n -H -u root /usr/local/libexec/hermes-email-boundary-verify
ExecStart=__BRIDGE_VENV__/bin/hermes-email-bridge poll --continuous
Restart=on-failure
RestartSec=30
UMask=0077
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=__STATE_DIR__ __LOG_DIR__ /var/lib/hermes-email-agent

[Install]
WantedBy=multi-user.target
